# provisioning/bootstrap/files/bootstrap-nat.service
[Unit]
Description=Bootstrap appliance NAT (PXE network -> uplink)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# BOOTSTRAP_IFACE is set during image build via envsubst
Environment=BOOTSTRAP_IFACE=__BOOTSTRAP_IFACE__
ExecStart=/usr/local/sbin/bootstrap-nat.sh
Restart=on-failure
RestartSec=2
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

